# relay-config.service
[Unit]
Description=Relay Control Web Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=gpio
WorkingDirectory=/home/pi/relay_control
Environment="PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/pi/relay_control/app.py
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/relay_control

# Resource limits
CPUQuota=50%
MemoryLimit=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay-control

[Install]
WantedBy=multi-user.target
