# relay-config.service
[Unit]
Description=Relay Control Web Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=tech
Group=gpio
WorkingDirectory=/home/tech/relay_control
Environment="PATH=/home/tech/relay_control/venv/bin"
ExecStart=/home/tech/relay_control/venv/bin/python3 /home/tech/relay_control/app.py
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# FIX: Add the line below to allow writing to the app's directory
ReadWritePaths=/home/tech/relay_control/

# Resource limits
CPUQuota=50%
MemoryLimit=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay-control

[Install]
WantedBy=multi-user.target
